# deploy/atlas.service
[Unit]
Description=Atlas Round-Robin Discord Bot
After=network.target postgresql.service

[Service]
Type=simple
# Change 'root' to your username if you cloned into home
User=root
# Adjust working directory to where you cloned the repo (e.g. /root/atlas)
WorkingDirectory=/opt/atlas-bot
ExecStart=/opt/atlas-bot/venv/bin/python main.py
Restart=on-failure
RestartSec=5
EnvironmentFile=/opt/atlas-bot/.env

[Install]
WantedBy=multi-user.target
